# backend/holido_supervisor.conf
[program:holido_api]
command = /work/Backend/holido/run.sh                  ; Command to start app
redirect_stderr=true            
stdout_logfile=/var/log/supervisor.txt